# root/etc/s6-overlay/s6-rc.d/svc-seatd/run
#!/usr/bin/with-contenv bash

set -xe

exec seatd -u $USER