# root/etc/s6-overlay/s6-rc.d/svc-sway/run
#!/usr/bin/with-contenv bash

set -xe

exec s6-setuidgid $USER \
  sway